# terraform/aws/cloud-resources/99-auto-lvm.rules
KERNEL=="nvme*n*", ACTION=="add", RUN+="/usr/local/bin/auto-lvm-add.sh /dev/%k", RUN+="/bin/systemctl start mnt-data-mount.service"
